# docker-trex/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:trex]
command=./t-rex-64 -i --cfg /etc/trex_cfg.yaml